# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

# ── API (FastAPI + uvicorn) ───────────────────────────────────────────────────
[program:api]
# PORT 由 entrypoint.sh 写入 /etc/supervisor/conf.d/env.conf 后通过 %(ENV_PORT)s 读取
# 若未注入则默认 8000
command=python3 -u -m uvicorn app.main:app --host 0.0.0.0 --port %(ENV_PORT)s
directory=/app
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB
redirect_stderr=false

# ── Worker-1 ──────────────────────────────────────────────────────────────────
[program:worker-1]
command=python3 -u worker/worker.py
directory=/app
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",STREAM_CONSUMER="worker-1"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/var/log/supervisor/worker-1.log
stderr_logfile=/var/log/supervisor/worker-1-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB

# ── Worker-2 ──────────────────────────────────────────────────────────────────
[program:worker-2]
command=python3 -u worker/worker.py
directory=/app
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",STREAM_CONSUMER="worker-2"
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/var/log/supervisor/worker-2.log
stderr_logfile=/var/log/supervisor/worker-2-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=5MB

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
